# Users/xuan/Library/LaunchAgents/com.bci.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.bci.gateway</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>cd "$HOME/Desktop/BCI Connect/nfc-gateway-service" && PORT=3002 CLOUD_API_URL="https://bci-connect.onrender.com" node server.js &gt; "$HOME/Desktop/BCI Connect/nfc-gateway-service/gateway.out.log" 2&gt; "$HOME/Desktop/BCI Connect/nfc-gateway-service/gateway.err.log"</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>${HOME}/Desktop/BCI Connect/nfc-gateway-service/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>${HOME}/Desktop/BCI Connect/nfc-gateway-service/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>PORT</key>
      <string>3002</string>
      <key>CLOUD_API_URL</key>
      <string>https://bci-connect.onrender.com</string>
    </dict>
  </dict>
</plist>